# rootfs/etc/services.d/mini_httpd/run
#!/usr/bin/with-contenv bash

mini_httpd -C /etc/mini_httpd/mini_httpd.conf -c "**.cgi|**.sh"  -l /dev/stdout -D